# docker/crontab
# GPS Tracker Cron Jobs
# Run as www-data user

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Fetch latest data every minute (actual frequency controlled by FETCH_FREQUENCY_MINUTES setting)
* * * * * www-data cd /var/www/html && /usr/local/bin/php fetch_data.php >> /var/log/tracker/fetch.log 2>&1

# Smart refetch every 30 minutes (detect missed/delayed data uploads)
# Uses SMART_REFETCH_DAYS from .env for days to check back
*/30 * * * * www-data cd /var/www/html && /usr/local/bin/php smart_refetch_v2.php >> /var/log/tracker/smart_refetch.log 2>&1

# Empty line required at end
